# system_resources/tors.service
[Unit]
After=network.service

[Service]
ExecStart=/root/Orb-Message-Board/system_resources/start_tors.sh

[Install]
WantedBy=default.target